# deploy/typingracebot.service
[Unit]
Description=Discord Typing Race Bot (Group 6)
After=network-online.target redis6.service
Wants=network-online.target redis6.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/bot
ExecStart=/usr/bin/java -jar /home/ec2-user/bot/discord-typing-race-1.0-SNAPSHOT-jar-with-dependencies.jar
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
