# deploy/typingracebot.service
[Unit]
Description=Discord Typing Race Bot (Group 6)
After=network-online.target redis6.service
Wants=network-online.target redis6.service

[Service]
Type=simple
User=ec2-user
<<<<<<< HEAD
WorkingDirectory=/Discord-bot-project-group-6
ExecStart=/usr/bin/java -jar /Discord-bot-project-group-6/discord-typing-race-1.0-SNAPSHOT-jar-with-dependencies.jar
=======
WorkingDirectory=/home/ec2-user/Discord-bot-project-group-6
ExecStart=/usr/bin/java -jar /home/ec2-user/Discord-bot-project-group-6/target/discord-typing-race-1.0-SNAPSHOT-jar-with-dependencies.jar
>>>>>>> e78373645e221ca80ceec9abe988ef237ff45bb1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
